# sslocal.conf
description "Service for shadowsocks client"
author "tsung shiu"

start on started ui
stop on stopping ui

script
    
    exec bash /etc/init.d/sslocal.sh

end script

